# launchctl-ruby.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>KeepAlive</key>
    <true/>
    <key>Label</key>
    <string><%= launch_agent_name %></string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/chrisroos/.rbenv/versions/2.2.1/bin/ruby</string>
        <string><%= script_path %></string>
    </array>
    <key>QueueDirectories</key>
    <array/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string><%= stderr_path %></string>
    <key>StandardOutPath</key>
    <string><%= stdout_path %></string>
    <key>WatchPaths</key>
    <array/>
    <key>WorkingDirectory</key>
    <string><%= current_directory %></string>
  </dict>
</plist>
